# 5.6/etc/services.d/php-fpm/finish
#!/bin/sh
# bring down container if php5-fpm stops
s6-svscanctl -t /var/run/s6/services
